# Automate Phishing via GoPhish _ N8N _ AI Prompts/.n8n/n8n.service
[Unit]
Description=N8N Workflow Automation
After=network.target

[Service]
Type=simple
User=$USER
Environment=N8N_HOST=0.0.0.0
Environment=N8N_PORT=5678
Environment=N8N_PROTOCOL=http
ExecStart=/usr/bin/n8n start
Restart=always

[Install]
WantedBy=multi-user.target
